# packages/network/openssh/system.d/sshd.service
[Unit]
Description=OpenSSH server daemon


[Service]
Restart=on-failure
ExecStartPre=/usr/bin/ssh-keygen -A
ExecStart=/usr/sbin/sshd 
ExecReload=/bin/kill -HUP $MAINPID
TimeoutStopSec=1s
RestartSec=2
StartLimitInterval=0

[Install]
WantedBy=multi-user.target
